# util/cron.d/fido-updater
MAILTO=root
*/5 * * * * root cd /root/Fido-Updater; git pull;
